# airsensor.service
[Unit]
Description=Pimoroni Enviro+ Air Quality Sensor Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/airsensor
ExecStart=/opt/airsensor/.venv/airsensor/bin/python /opt/airsensor/src/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
