# airsensor.service
[Unit]
Description=Air Quality Sensor Monitor
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=airsensor
Group=airsensor
WorkingDirectory=/opt/airsensor

ExecStart=/opt/airsensor/.venv/bin/python3 /opt/airsensor/src/main.py

# Restart policy
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=airsensor

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/airsensor/logs

# Allow access to hardware
SupplementaryGroups=gpio i2c spi dialout

[Install]
WantedBy=multi-user.target
